# seo-pulse-web.service
[Unit]
Description=SEO Pulse Next.js Web App
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/seo-pulse
ExecStart=/usr/bin/npm start
Restart=always
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
